# overlay/etc/systemd/system/serial-getty@ttyS0.service.d/bounded-restart.conf
# Bound the agetty respawn loop on hardware where /dev/ttyS0 exists as a
# chardev (8250 driver registered it from platform/ACPI tables) but has no
# working UART backing. tcgetattr returns EIO, agetty exits, systemd
# respawns. The stock template sets Restart=always RestartSec=0, which on
# such boxes fills the journal with one failure every ~10s forever.
#
# Cap to 3 starts per 120s; after that the unit goes to `failed` and stays
# there for the rest of the boot. Operators with a real serial port aren't
# affected — the first start succeeds and we never hit the limit.

[Unit]
StartLimitIntervalSec=120
StartLimitBurst=3

[Service]
Restart=on-failure
RestartSec=10
