# overlay/etc/systemd/system/vestick-overlay-resize.service
[Unit]
Description=VEstick first-boot resize of f2fs overlay
DefaultDependencies=no
After=local-fs.target
Before=basic.target sysinit.target
ConditionPathExists=!/var/lib/vestick/.overlay-resized

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/vestick-overlay-resize

[Install]
WantedBy=sysinit.target
