# overlay/etc/systemd/system/veyage-firstboot.service
[Unit]
Description=VEyage first-boot identity setup (interactive: hostname + root password)
After=veyage-state.service
Before=veyage-network-init.service networking.service
RequiresMountsFor=/mnt/state

[Service]
Type=oneshot
RemainAfterExit=yes
StandardInput=tty
StandardOutput=tty
StandardError=tty
TTYPath=/dev/console
TTYReset=yes
TTYVHangup=yes
TimeoutStartSec=infinity
ExecStart=/usr/local/sbin/veyage-firstboot

[Install]
WantedBy=networking.service
