# overlay/etc/systemd/system/veyage-overlay-resize.service
[Unit]
Description=VEyage first-boot resize of f2fs overlay
DefaultDependencies=no
After=local-fs.target
Before=basic.target sysinit.target
ConditionPathExists=!/var/lib/veyage/.overlay-resized

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/veyage-overlay-resize

[Install]
WantedBy=sysinit.target
