# overlay/etc/systemd/system/veyage-firstboot.service
[Unit]
Description=VEyage first-boot identity setup (interactive: SSH keys + hostname + root password)
Before=veyage-network-init.service networking.service getty.target
ConditionPathExists=!/var/lib/veyage/.firstboot-done

[Service]
Type=oneshot
RemainAfterExit=yes
StandardInput=tty
StandardOutput=tty
StandardError=tty
TTYPath=/dev/console
TTYReset=yes
TTYVHangup=yes
TimeoutStartSec=infinity
ExecStart=/usr/local/sbin/veyage-firstboot

[Install]
WantedBy=networking.service
