# overlay/etc/systemd/system/veyage-resize-root.service
[Unit]
Description=VEyage first-boot rootfs resize (grow to fill device)
ConditionPathExists=!/var/lib/veyage/resized
After=local-fs.target
Before=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/veyage-resize-root

[Install]
WantedBy=multi-user.target
