# overlay/etc/systemd/system/veyage-firstboot.service
[Unit]
Description=VEyage first-boot identity setup (systemd-firstboot in an rw window)
ConditionFirstBoot=yes
DefaultDependencies=no
After=systemd-remount-fs.service systemd-tmpfiles-setup.service systemd-vconsole-setup.service
Before=sysinit.target shutdown.target
Conflicts=shutdown.target
Wants=first-boot-complete.target
Before=first-boot-complete.target

[Service]
Type=oneshot
RemainAfterExit=yes
StandardInput=tty
StandardOutput=tty
StandardError=tty
TTYPath=/dev/console
TTYReset=yes
TTYVHangup=yes
TimeoutStartSec=infinity
# / is mounted ro per fstab. Briefly flip it rw, run systemd-firstboot
# (which prompts at /dev/console for hostname + root password and writes
# /etc/hostname + /etc/shadow), then flip back. No bind-mount inode
# gymnastics — it's just normal ext4 writes now.
ExecStartPre=/bin/mount -o remount,rw /
ExecStart=/usr/bin/systemd-firstboot --prompt-hostname --prompt-root-password
ExecStartPost=/bin/sync
ExecStartPost=/bin/mount -o remount,ro /

# Pick up basic fields from systemd credentials (so an unattended deploy
# can preseed root password / hostname via SMBIOS or kernel cmdline).
ImportCredential=passwd.hashed-password.root
ImportCredential=passwd.plaintext-password.root
ImportCredential=firstboot.*

[Install]
WantedBy=sysinit.target
