# overlay/etc/systemd/system/pve-cluster.service.d/veyage-wait-for-hosts.conf
# pve-cluster's pmxcfs binary calls gethostbyname() at start; if the
# local hostname maps only to 127.0.0.1, it logs "missing entry in
# /etc/hosts or DNS" and exits 255. veyage-hosts-init.service writes a
# real-IP hosts entry to /run/hosts and bind-mounts it over /etc/hosts;
# this drop-in just enforces ordering so pmxcfs sees the entry.
[Unit]
After=veyage-hosts-init.service
Wants=veyage-hosts-init.service
