# overlay/etc/systemd/system/pve-firewall.service.d/veyage-no-update-alternatives.conf
# pve-firewall.service ships ExecStartPre lines that run
# `update-alternatives --set ebtables/iptables/ip6tables /usr/sbin/*-legacy`
# at every boot. update-alternatives --set always writes
# /etc/alternatives/<name>.dpkg-tmp before atomic-renaming, which fails on
# our ro /etc. The default state for these alternatives is set to legacy at
# build time (build.sh cleanup_image), so the runtime --set is redundant.
#
# Empty ExecStartPre= clears the inherited list.
[Service]
ExecStartPre=
