# overlay/etc/systemd/system/pveproxy.service.d/veyage-no-pveupdate.conf
# pveproxy.service ships an ExecStartPost that conditionally fires
# /usr/bin/pveupdate (gated on /var/log/pveam.log not existing — true on
# every boot since /var/log is tmpfs). pveupdate refreshes apt sources
# and writes /etc/apt/apt.conf.d/76pveproxy via atomic rename, which
# fails on ro /etc with "unable to open .tmp.<PID> - Read-only file
# system" each boot.
#
# On this design apt updates are operator-driven inside a remountrw
# window (the operator's `apt update && apt upgrade` will reach the
# necessary state), so the auto-fire serves no purpose. Empty
# ExecStartPost= clears the inherited list.
[Service]
ExecStartPost=
