# overlay/etc/systemd/system/ssh.service.d/veyage-keygen.conf
# Generate SSH host keys on first start if missing. ssh-keygen -A is
# idempotent — it only creates the keys that don't exist, so this is a
# no-op after the first boot. /etc/ssh sits on the persistent f2fs
# overlay, so keys created here survive reboots.
#
# Replaces the build-time ssh-keygen -A in build.sh — without baked-in
# keys, every dd'd device gets unique keys generated against its own
# /dev/random.
[Service]
ExecStartPre=/usr/bin/ssh-keygen -A
